# applications/luci-app-multi-user/root/etc/uci-defaults/09_users
#!/bin/sh
touch /etc/config/users
echo "config user 'new'" > /etc/config/users
exit 0
